# startup_scripts/raspbian/eyeballs.service
# Automatically generated by systemd-sysv-generator

[Install]
WantedBy=multi-user.target

[Unit]
SourcePath=/etc/init.d/eyeballs
Description=SYSV: eyeballs daemon

[Service]
Type=forking
Restart=no
TimeoutSec=5min
IgnoreSIGPIPE=no
KillMode=process
GuessMainPID=no
RemainAfterExit=yes
ExecStart=/etc/init.d/eyeballs start
ExecStop=/etc/init.d/eyeballs stop
